# devops/lib/dpkg/info/virtualbox-guest-utils.list
/.
/etc
/etc/init.d
/etc/init.d/virtualbox-guest-utils
/lib
/lib/security
/lib/security/pam_vbox.so
/lib/systemd
/lib/systemd/system
/lib/systemd/system/virtualbox-guest-utils.service
/usr
/usr/sbin
/usr/sbin/VBoxService
/usr/sbin/VBoxControl
/usr/share
/usr/share/doc
/usr/share/doc/virtualbox-guest-utils
/usr/share/doc/virtualbox-guest-utils/changelog.Debian.gz
/usr/share/doc/virtualbox-guest-utils/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/mount.vboxsf.8.gz
/usr/share/man/man8/VBoxControl.8.gz
/usr/share/man/man8/VBoxService.8.gz
/sbin
/sbin/mount.vboxsf
